# rootfs/build/config/supervisord_x.conf
[supervisord]
logfile = /tmp/supervisord_x.log
logfile_maxbytes = 0
pidfile = /tmp/supervisord_x.pid

[program:openbox]
command = openbox
autostart = true
autorestart = true
startsecs = 1
priority = 10
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/fd/2
stderr_logfile_maxbytes = 0

[program:socat]
command = socat TCP-LISTEN:9222,fork,bind=0.0.0.0,reuseaddr TCP:127.0.0.1:19222
autostart = true
autorestart = true
startsecs = 1
priority = 20
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/fd/2
stderr_logfile_maxbytes = 0

[program:chrome]
command = /build/bin/start_chrome
autostart = true
autorestart = true
startsecs = 1
priority = 30
stopsignal = TERM
stopasgroup = true
killasgroup = true
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/fd/2
stderr_logfile_maxbytes = 0
