# rootfs/build/services/vnc.conf
[program:vnc]

command = /build/bin/start_vnc
autostart = true
autorestart = true
