# rootfs/.init/services/cdp.conf
[program:cdp]

command = /build/bin/start_cdp
autostart = true
autorestart = true
stopsignal = TERM
stopasgroup = true
killasgroup = true
