# rootfs/.init/services/xvfb.conf
[program:xvfb]

command = /build/bin/start_xvfb
autostart = true
autorestart = true
priority = 5
